# systemd/openkoutsi-frontend@.service
[Unit]
Description=openkoutsi frontend
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/projects/openkoutsi/frontend
Environment=PORT=3000
# Adjust the npm path if you switch Node versions via nvm.
ExecStart=%h/.nvm/versions/node/v25.2.1/bin/npm start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
